# etc/cron.d/owncloud
*/15 *	* * *	fcgi-php	cd /var/www/owncloud && php ./cron.php
